# root/etc/services.d/mstream/run
#!/usr/bin/with-contenv bash

cd /opt/mstream || exit

exec \
	s6-setuidgid abc /opt/mstream/mstream -d /config/db -I /config/album-art -j /config/config.json -m /music -u ${USER} -x ${PASSWORD} 