# root/etc/services.d/mstream/run
#!/usr/bin/with-contenv bash

cd /opt/mstream || exit

exec \
    s6-setuidgid abc mstream -j /config/config.json